# Winstride-Web/src/shared/detection/sigma/rules/powershell/ccdc_ps_scheduled_task.yml
title: Scheduled Task Created via PowerShell
id: b0000002-0000-0000-0000-000000000021
name: ps_scheduled_task
status: stable
description: A PowerShell script creates or registers a scheduled task — common persistence mechanism.
author: WinStride CCDC
date: 2025-01-01
tags:
    - attack.persistence
    - attack.t1053.005
logsource:
    category: ps_script
    product: windows
detection:
    selection:
        ScriptBlockText|contains:
            - 'New-ScheduledTask'
            - 'Register-ScheduledTask'
            - 'schtasks'
    condition: selection
level: medium
